# build_files/config/includes.chroot/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/conky
Hidden=false
X-GNOME-Autostart-enabled=true
Name=conky
