# build_files/config/includes.chroot/etc/skel/.config/autostart/readme.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/mousepad /home/thor/Desktop/README.txt
Hidden=false
X-GNOME-Autostart-enabled=true
Name=readme

